# sk5_auto.sh
#!/bin/bash
#chkconfig: 2345 80 90
    #description:auto_run
/etc/init.d/sockd start
