# sk5_auto.sh
#!/bin/bash
#chkconfig: 2345 80 90
    #description:sk5_auto
/etc/init.d/sockd start
